# systemd/user/git-sync.timer
[Unit]
Description=Run git sync every 5 minutes

[Timer]
OnCalendar=*:0/2
Persistent=true

[Install]
WantedBy=timers.target
