# systemd/user/git-sync.service
[Unit]
Description=Sync all git repos
After=network.target

[Service]
Type=oneshot
ExecStart=/home/catalin/.local/bin/git-sync-all.sh
