# systemd/user/git-sync@.service
[Unit]
Description=Auto-sync Git repo at %f
After=network.target

[Service]
Type=oneshot
# %f provides the unescaped path automatically in newer systemd versions
WorkingDirectory=%f
ExecStart=/usr/bin/bash -c 'git add . && (git diff-index --quiet HEAD || git commit -m "Auto-sync [$(hostname)] - $(date +%%Y-%%m-%%d\ %%H:%%M:%%S)") && git pull --rebase --autostash && git push'
[Install]
WantedBy=default.target
